# etc/profile.d/inputrc.sh
# Overwrite inputrc with our own customizations
export INPUTRC=/usr/share/solydxk/default-settings/inputrc
